# etc/profile.d/pulse.sh
pulseaudio -D --high-priority=true --realtime=true --exit-idle-time=-1 --no-cpu-limit